# ssadan.service
[Unit]
Description=SSADAN Meal Notification Bot
After=network.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/path/to/ssadan
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /path/to/ssadan/main.py run
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ssadan

[Install]
WantedBy=multi-user.target
